# etc/supervisord.d/postfix.ini
[program:postfix]
process_name = master
command=/usr/sbin/postfix -c /etc/postfix start

priority=100
startsecs=0
; autostart=false
autorestart=false