# etc/supervisord.d/rsyslog.ini
[program:rsyslog]
command=/usr/sbin/rsyslogd -n

; startsecs=2

priority=200
; autostart=true
; autorestart=false